# installer/macos/com.workflow-recorder.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
LaunchAgent template for Workflow Recorder.

The postinstall script substitutes the tokens below before writing this
file to ~/Library/LaunchAgents/com.workflow-recorder.plist:

  __HOME__     absolute path of the user's home directory
  __BUNDLE__   /Applications/WorkflowRecorder.app
  __CONFIG__   ~/Library/Application Support/WorkflowRecorder/model_config.json

Runs as a user-level agent — the recorder needs access to the logged-in
user's desktop, so a system LaunchDaemon would see the login window
instead of the user's screen.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.workflow-recorder</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BUNDLE__/Contents/MacOS/workflow-recorder</string>
        <string>-c</string>
        <string>__CONFIG__</string>
    </array>

    <!-- Auto-start on user login and relaunch on crash. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <!-- Throttle relaunches to at most once every 30s so a broken
         config doesn't pin a core. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/WorkflowRecorder/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/WorkflowRecorder/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>__HOME__/Library/Application Support/WorkflowRecorder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin</string>
    </dict>

    <!-- ProcessType: Interactive so the agent gets normal QoS priority
         for screen capture (vs. Background which throttles CPU). -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
